# systemd/carolins-kasse-backup.timer
[Unit]
Description=Daily Carolin's Kasse database backup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
