# systemd/carolins-install.service
[Unit]
Description=Install Carolin's Kasse after first boot
After=network-online.target time-sync.target
Wants=network-online.target time-sync.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/carolins-kasse/install.env
ExecStart=/usr/local/sbin/carolins-pi-bootstrap.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
